# bsrouter/bsrouter.service
[Unit]
Description=bsrouter service
After=network.target

[Service]
User=bsrouter
Restart=always
WorkingDirectory=~
ExecStart=/usr/local/bin/bsrouter

[Install]
WantedBy=multi-user.target